# etc/profile.d/zpath.sh
# -*- mode: sh; sh-basic-offset: 2; indent-tabs-mode: nil; -*-
# vim:set ft=sh et sw=2 ts=2:
#
# zpath.sh - cleanup apath.sh functions

unset -f add_root_path add_path remove_path
